# etc/udev/rules.d/20-udev-fb.rules
KERNEL=="fb0", ATTR{name}=="fb_ili9340", ACTION=="add", RUN+="/etc/cloudshell/cloudshell.sh"
